# examples/bye-axum/axum.socket
[Unit]
Description=Axum Socket

[Socket]
ListenStream=0.0.0.0:8080
Accept=no

[Install]
WantedBy=sockets.target
